# build/docker/crontab


0 1 * * * curl http://web.spy/update-user-process?key=fgjegrfiugsddkfjhsd >> /var/log/cron.log 2>&1
